# airmx-addon/rootfs/etc/services.d/i-airmx-cn/run
#!/usr/bin/with-contenv bashio
exec /i-airmx-cn.py
