# airmx-addon/rootfs/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv bashio
exec mosquitto -c /etc/mosquitto/mosquitto.conf
